# hawkbook/_root/etc/systemd/system/hawkroot.service
[Unit]
Description=Startup stuff for Hawkbook

[Service]
User=root
ExecStart=/home/hawkens/startup/hawkroot.sh

[Install]
WantedBy=multi-user.target
